# deploy/llcprime-api.service
# systemd unit for the LGU PRIME-HRM backend (Node/Express, port 5010)
# ---------------------------------------------------------------------------
# Install:
#   sudo cp deploy/llcprime-api.service /etc/systemd/system/llcprime-api.service
#   # edit User/WorkingDirectory/paths below to match the server
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now llcprime-api
#   sudo systemctl status llcprime-api
#   journalctl -u llcprime-api -f      # live logs
#
# Use this OR PM2 (deploy/ecosystem.config.cjs) — not both.
# ---------------------------------------------------------------------------
[Unit]
Description=LGU PRIME-HRM API (llcprime)
After=network.target mysql.service

[Service]
Type=simple
# Match an existing non-root deploy user on the box (e.g. www-data or a dedicated user)
User=www-data
Group=www-data
WorkingDirectory=/var/www/llcprime/server
# Reads server/.env for DATABASE_URL, JWT secrets, PORT, CORS_ORIGIN
EnvironmentFile=/var/www/llcprime/server/.env
# Absolute path to node — check with `which node`; adjust if using nvm
ExecStart=/usr/bin/node dist/app.js
Restart=on-failure
RestartSec=5
# Hardening
NoNewPrivileges=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
